# systemd/archive-snapshot.timer
[Unit]
Description=Monthly Archive Billing Snapshot Timer
Requires=archive-snapshot.service

[Timer]
# Run on the 1st of each month at 2:00 AM
OnCalendar=*-*-01 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
